# docs/botjohn-saturday-maintenance.service
[Unit]
Description=BotJohn — Saturday research maintenance (Sat 16:00 ET)
After=network.target postgresql.service openclaw-saturday-brain.service
Wants=openclaw-saturday-brain.service

[Service]
Type=oneshot
User=claudebot
Group=claudebot
EnvironmentFile=/root/openclaw/.env
WorkingDirectory=/root/openclaw
ExecStart=/usr/bin/node /root/openclaw/src/agent/run_maintenance.js --mode saturday
StandardOutput=journal
StandardError=journal
TimeoutStartSec=3600

[Install]
WantedBy=multi-user.target
