# docs/mastermind-weekly.timer
[Unit]
Description=Run MastermindJohn (strategy-stack mode) every Friday at 20:00 America/New_York

[Timer]
OnCalendar=Fri *-*-* 20:00:00 America/New_York
Persistent=true
Unit=openclaw-mastermind-weekly.service

[Install]
WantedBy=timers.target
