# docs/mastermind-weekly.service
[Unit]
Description=MastermindJohn — weekly strategy-stack review (Friday)
After=network.target postgresql.service
PartOf=johnbot.service

[Service]
Type=oneshot
User=claudebot
Group=claudebot
EnvironmentFile=/root/openclaw/.env
WorkingDirectory=/root/openclaw
ExecStart=/usr/bin/node /root/openclaw/src/agent/curators/run_mastermind.js --mode strategy-stack
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
