# docs/verdict-refresh.timer
[Unit]
Description=Run verdict-refresh wave every 4 hours off-market

[Timer]
# 02/06/18/22 UTC = 22/02/14/18 ET (off-market on US-equity windows). Adjust
# if you want to align to your trading hours.
OnCalendar=*-*-* 02,06,18,22:00:00
Persistent=true
Unit=openclaw-verdict-refresh.service

[Install]
WantedBy=timers.target
